# ssl/reminder.cron
0 0 * * * php /var/www/ssl/reminder.php > /dev/null
